# erica/cron.d/delete_success_fail_entities
SHELL=/bin/bash
*/5 * * * * root cd /app && /usr/local/bin/pipenv run python erica/shared/sqlalchemy/cron/update_entities_utils.py delete-success-fail-entities &> /app/cronjob_success_fail_output
