# erica/cron.d/set_not_processed_entities_to_failed
SHELL=/bin/bash
* * * * * root cd /app && /usr/local/bin/pipenv run python erica/erica_shared/sqlalchemy/cron/update_entities_utils.py set-not-processed-entities-to-failed &> /app/cronjob_not_processed_output
